# supervisor/supervisord.conf
[supervisord]
logfile = /dev/stdout
logfile_maxbytes = 0
loglevel = info
nodaemon = true

[unix_http_server]
file = /tmp/supervisor.sock
chmod = 0777

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock
prompt = supervisor

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:db]
command = /usr/local/bin/mysqld-startup.sh
directory = /
user = mysql
autostart = false
autorestart = unexpected
startsecs = 10
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0

[program:app]
command = /usr/local/bin/app-startup.sh
directory = /AWSBPM/bin
user = root
autostart = false
autorestart = unexpected
startsecs = 10
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0

[program:web]
command = /usr/local/bin/web-startup.sh
directory = /AWSBPM/bin
user = root
autostart = false
autorestart = unexpected
startsecs = 10
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0