# supervisor/supervisord.ini
[supervisord]
logfile = /dev/stdout
logfile_maxbytes = 0
loglevel = info
nodaemon = true

[program:db]
command=/usr/local/bin/mysql-startup.sh
directory=/
user=mysql
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0

[program:app]
command=/usr/local/bin/app-startup.sh
directory=/AWSBPM/bin
user=root
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0

[program:web]
command=/usr/local/bin/web-startup.sh
directory=/AWSBPM/bin
user=root
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0